# deploy/systemd/nginx.service.d/bwb-fiscal-open-boot-recovery.conf
# Drop-in: nginx.service must not start unless boot recovery succeeded.
# Armed reboot with failed recovery → nginx.service start is refused (Requires=).
# Confirmed: boot recovery no-ops successfully → nginx may start open.

[Unit]
Requires=bwb-fiscal-nginx-open-boot-recovery.service
After=bwb-fiscal-nginx-open-boot-recovery.service
