# deploy/systemd/bwb-fiscal-nginx-open-boot-recovery.service
# Runs BEFORE nginx.service so an armed reboot cannot accept public traffic
# while still open. Restores deny-all on disk and validates nginx -t only
# (no reload/curl — nginx is not up yet). confirmed leaves open config intact.

[Unit]
Description=BWB fiscal Nginx open boot recovery (fail-closed if armed)
Documentation=file:///opt/bwb-modulo-fiscal/current/nginx/README.md
DefaultDependencies=no
Requires=local-fs.target
After=local-fs.target
Before=nginx.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bwb-fiscal-deploy-helper nginx-open-boot-recovery
RemainAfterExit=yes

[Install]
WantedBy=nginx.service
WantedBy=multi-user.target
